# systemd/fanctl-boot-apply.service
[Unit]
Description=Re-apply last manual fan mode after boot
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fanctl restore
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
